# scripts/directorio.conf
# myservice - myservice job file

description "Un servicio para inciar el aplicacion de directorio con Node.js y el modulo forever"
author "Andrew Waterman <awaterma@awaterma.net>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Run before process
pre-start script
    [ -d /var/run/directorio ] || mkdir -p /var/run/directorio
end script

# Start the process
exec forever start /opt/directorio/src/app.js
  
